# backend/scraper/crontab
* * * * * flask scrape-v2 > /backend/scraper/scrape.log 2>&1
